# kbp-webalizer/files/cron.daily/webalizer
#!/bin/sh

su www-data -c /usr/local/bin/webalizer-multi

# Exit with the last exit code
exit $?
